# crontab
* * * * * root echo "1 minute later ..." &>> /var/log/cron.log
# asdasd
